# _docker/api/stage.supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=/usr/local/bin/uwsgi /app/_docker/api/stage.uwsgi.ini

