# _docker/api/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=/usr/local/bin/uwsgi /var/www/quicksilver_api/_docker/api/uwsgi.ini

